# nectar.service
[Unit]
Description=Nectar Cardano Blockchain Indexer
After=network.target tidb.service

[Service]
Type=simple
User=cardano
Group=cardano
WorkingDirectory=/opt/nectar
Environment="DASHBOARD_TYPE=both"
Environment="WEB_PORT=8080"
EnvironmentFile=/opt/nectar/.env
ExecStart=/opt/nectar/nectar
Restart=on-failure
RestartSec=10

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nectar/logs

[Install]
WantedBy=multi-user.target